# roles/docker-compose/files/docker-compose-reload.service
[Unit]
Description=Refresh images and update containers

[Service]
Type=oneshot

ExecStart=/usr/bin/systemctl reload docker-compose.service
